# crontab
# Roda a cada 2 horas carregando as variáveis de ambiente
0 */2 * * * . /etc/environment; cd /app && /usr/local/bin/python -m app.main >> /var/log/cron.log 2>&1
